# daemonize-packaging/debian/postinst
#!/bin/sh
# postinst script for aptly
#
# see: dh_installdeb(1)

set -e

# summary of how this script can be called:
#        * <postinst> `configure' <most-recently-configured-version>
#        * <old-postinst> `abort-upgrade' <new version>
#        * <conflictor's-postinst> `abort-remove' `in-favour' <package>
#          <new-version>
#        * <postinst> `abort-remove'
#        * <deconfigured's-postinst> `abort-deconfigure' `in-favour'
#          <failed-install-package> <version> `removing'
#          <conflicting-package> <version>
# for details, see https://www.debian.org/doc/debian-policy/ or
# the debian-policy package


action="$1"

set_permissions()
{
  if ! getent passwd aptlyuser >/dev/null; then
    adduser --system --group --no-create-home --quiet aptlyuser
  fi
  chown -R aptlyuser:aptlyuser /var/aptly
}

set_boot()
{
  if [ -x /etc/init.d/aptly ]; then
    update-rc.d aptly defaults >/dev/null || exit $?
  fi
}


case "$action" in
    configure)
      set_permissions
      set_boot
    ;;

    abort-upgrade|abort-remove|abort-deconfigure)
    ;;

    *)
        echo "postinst called with unknown argument \`$action'" >&2
        exit 1
    ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0
